# vllm-bootc/etc/systemd/system/rhoim-vllm.service
[Unit]
Description=RHOIM vLLM Model Serving Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Run as the rhoim service user created via sysusers.d/rhoim.conf
User=rhoim
Group=rhoim

# Load configuration from /etc/sysconfig/rhoim (MODEL_ID, MODEL_PATH, VLLM_DEVICE_TYPE, etc.)
EnvironmentFile=-/etc/sysconfig/rhoim

# Minimal defaults; can be overridden by the env file or systemd overrides
Environment=VLLM_MODEL=TinyLlama/TinyLlama-1.1B-Chat-v1.0
Environment=HOST=0.0.0.0
Environment=PORT=8000
Environment=DTYPE=float32
Environment=VLLM_EXTRA_ARGS=

# Do NOT set CPU-only envs here; let the script choose based on VLLM_DEVICE_TYPE
ExecStart=/usr/local/bin/initializer-entrypoint.sh

StandardOutput=journal+console
StandardError=journal+console
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target