# vllm-bootc/etc/systemd/system/rhoim-vllm.service
[Unit]
Description=RHOIM vLLM Model Serving Service (Podman Container)
After=network-online.target podman.socket nvidia-container-setup.service
Wants=network-online.target
Requires=podman.socket

[Service]
Type=simple

# Run the wrapper script that handles environment variables and podman execution
ExecStart=/usr/local/bin/run-vllm-container.sh

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
