# deploy/bootc/rhoim-vllm.service
[Unit]
Description=RHOIM vLLM Model Serving Service
After=network-online.target

[Service]
Type=simple
# Load environment variables from config file
EnvironmentFile=-/etc/sysconfig/rhoim
# ExecStart calls the initializer script
ExecStart=/usr/local/bin/initializer-entrypoint.sh
# Redirect output to console for debugging
StandardOutput=journal+console
StandardError=journal+console
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
