# deploy/bootc/rhoim-vllm.service
[Unit]
Description=RHOIM vLLM Model Serving Service
After=network-online.target

[Service]
Type=simple
# NEW: ExecStart calls the initializer script
ExecStart=/usr/local/bin/initializer-entrypoint.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
